# Code And Services/catalis-led.service
[Unit]
Description=Control LEDs 

[Service]
Type=simple
ExecStart=/root/new.sh

[Install]
WantedBy=multi-user.target
